# files/etc/crontabs/root
0 2 * * * /usr/share/cmzj/openwrt-k_tool.sh update rules
